# scripts/readahead.runlevel-default.list
//usr/lib/gcc/i686-pc-linux-gnu/3.4.3/libstdc++.so.6
/bin/domainname
/bin/echo
/bin/egrep
/bin/fuser
/bin/grep
/bin/kill
/bin/login
/bin/more
/bin/mount
/bin/netstat
/bin/ping
/bin/ps
/bin/red
/bin/sort
/bin/umount
/etc/asound.state
/etc/autofs/auto.master
/etc/conf.d/domainname
/etc/conf.d/local.start
/etc/conf.d/local.stop
/etc/conf.d/ntp-client
/etc/dnsdomainname
/etc/exports
/etc/fam.conf
/etc/gentoo-release
/etc/init.d/alsasound
/etc/init.d/autofs
/etc/init.d/cupsd
/etc/init.d/distccd
/etc/init.d/domainname
/etc/init.d/esound
/etc/init.d/famd
/etc/init.d/functions.sh
/etc/init.d/hotplug
/etc/init.d/local
/etc/init.d/net.eth0
/etc/init.d/netmount
/etc/init.d/nfs
/etc/init.d/nfsmount
/etc/init.d/nscd
/etc/init.d/ntp-client
/etc/init.d/ntpd
/etc/init.d/slpd
/etc/init.d/sshd
/etc/init.d/svscan
/etc/init.d/sysklogd
/etc/init.d/vixie-cron
/etc/make.conf
/etc/nsswitch.conf
/etc/ntp.conf
/etc/portage/package.mask
/etc/profile
/etc/resolv.conf
/etc/ssh/ssh_host_dsa_key
/etc/ssh/ssh_host_key
/etc/ssh/ssh_host_rsa_key
/etc/ssh/sshd_config
/lib/cpp
/lib/ld-linux.so.2
/lib/libacl.so.1
/lib/libattr.so.1
/lib/libcrypt.so.1
/lib/libdl.so.2
/lib/libncurses.so.5
/lib/libnsl.so.1
/lib/libpam.so.0
/lib/libpam_misc.so.0
/lib/libproc-3.2.5.so
/lib/libresolv.so.2
/lib/libutil.so.1
/lib/libwrap.so.0
/lib/libz.so.1
/lib/tls/libc.so.6
/lib/tls/libm.so.6
/lib/tls/libpthread.so.0
/lib/tls/librt.so.1
/sbin/lsmod
/sbin/modprobe
/sbin/rmmod
/sbin/rpc.statd
/sbin/start-stop-daemon
/tmp/debug
/usr/bin/[
/usr/bin/access
/usr/bin/aconnect
/usr/bin/as
/usr/bin/awk
/usr/bin/c++
/usr/bin/cat
/usr/bin/cc
/usr/bin/ccache
/usr/bin/chmod
/usr/bin/chown
/usr/bin/client
/usr/bin/co
/usr/bin/cp
/usr/bin/cpp
/usr/bin/current
/usr/bin/cvs
/usr/bin/date
/usr/bin/diff
/usr/bin/digest
/usr/bin/dir
/usr/bin/dirname
/usr/bin/disable
/usr/bin/display
/usr/bin/distcc
/usr/bin/distccd
/usr/bin/ebuild
/usr/bin/echo
/usr/bin/emerge
/usr/bin/enable
/usr/bin/env
/usr/bin/esd
/usr/bin/false
/usr/bin/file
/usr/bin/find
/usr/bin/free
/usr/bin/fuser
/usr/bin/g++
/usr/bin/g77
/usr/bin/gawk
/usr/bin/gcc-config
/usr/bin/gcj
/usr/bin/getopt
/usr/bin/gpg
/usr/bin/host
/usr/bin/hostname
/usr/bin/id
/usr/bin/import
/usr/bin/indent
/usr/bin/info
/usr/bin/install
/usr/bin/killall
/usr/bin/lambda
/usr/bin/line
/usr/bin/ln
/usr/bin/locate
/usr/bin/logger
/usr/bin/man
/usr/bin/merge
/usr/bin/mkdir
/usr/bin/mv
/usr/bin/net
/usr/bin/newer
/usr/bin/nice
/usr/bin/passwd
/usr/bin/portageq
/usr/bin/printenv
/usr/bin/prune
/usr/bin/quota
/usr/bin/rev
/usr/bin/rm
/usr/bin/rsync
/usr/bin/sed
/usr/bin/server
/usr/bin/size
/usr/bin/sleep
/usr/bin/sort
/usr/bin/ssh
/usr/bin/stat
/usr/bin/strings
/usr/bin/svc
/usr/bin/svscan
/usr/bin/sync
/usr/bin/test
/usr/bin/time
/usr/bin/touch
/usr/bin/tr
/usr/bin/true
/usr/bin/uname
/usr/bin/uniq
/usr/bin/w
/usr/bin/which
/usr/bin/yes
/usr/lib/MozillaFirefox/libnspr4.so
/usr/lib/MozillaFirefox/libnss3.so
/usr/lib/MozillaFirefox/libplc4.so
/usr/lib/MozillaFirefox/libplds4.so
/usr/lib/MozillaFirefox/libsmime3.so
/usr/lib/MozillaFirefox/libsoftokn3.so
/usr/lib/MozillaFirefox/libssl3.so
/usr/lib/autofs/autofs-ldap-auto-master
/usr/lib/binutils/i686-pc-linux-gnu/2.15.92.0.2/libbfd-2.15.92.0.2.so
/usr/lib/gcc-lib/i686-pc-linux-gnu/3.3.4/libstdc++.so.5
/usr/lib/gcc/i686-pc-linux-gnu/3.4.3/libgcc_s.so.1
/usr/lib/graphviz/libagraph.so.0
/usr/lib/graphviz/libcdt.so.0
/usr/lib/libICE.so.6
/usr/lib/libMagick.so.6
/usr/lib/libSM.so.6
/usr/lib/libX11.so.6
/usr/lib/libXext.so.6
/usr/lib/libXi.so.6
/usr/lib/libXt.so.6
/usr/lib/libasound.so.2
/usr/lib/libaudiofile.so.0
/usr/lib/libbz2.so.1.0
/usr/lib/libcrypto.so.0.9.7
/usr/lib/libcups.so.2
/usr/lib/libesd.so.0
/usr/lib/libfreetype.so.6
/usr/lib/libgdk-1.2.so.0
/usr/lib/libglib-1.2.so.0
/usr/lib/libgmodule-1.2.so.0
/usr/lib/libgthread-1.2.so.0
/usr/lib/libgtk-1.2.so.0
/usr/lib/libjpeg.so.62
/usr/lib/liblber-2.2.so.7
/usr/lib/libldap-2.2.so.7
/usr/lib/libltdl.so.3
/usr/lib/libmagic.so.1
/usr/lib/libpopt.so.0
/usr/lib/libslp.so.1
/usr/lib/libssl.so.0.9.7
/usr/lib/libtiff.so.3
/usr/lib/libxmms.so.1
/usr/lib/misc/gcc-config
/usr/lib/portage/bin/emerge
/usr/local/bin/h
/usr/local/bin/n
/usr/local/bin/stats
/usr/sbin/alsactl
/usr/sbin/automount
/usr/sbin/cron
/usr/sbin/cupsd
/usr/sbin/exportfs
/usr/sbin/famd
/usr/sbin/nscd
/usr/sbin/nscd_nischeck
/usr/sbin/ntpd
/usr/sbin/rpc.mountd
/usr/sbin/rpc.nfsd
/usr/sbin/rpc.rquotad
/usr/sbin/slpd
/usr/sbin/sshd
/usr/share/karamba/bin/open
/usr/share/ntp/ntp.conf
