# scripts/readahead.runlevel-boot.list
//usr/lib/gcc/i686-pc-linux-gnu/3.4.3/libgcc_s.so.1
//usr/lib/gcc/i686-pc-linux-gnu/3.4.3/libstdc++.so.6
/bin/dmesg
/bin/echo
/bin/egrep
/bin/gawk
/bin/grep
/bin/hostname
/bin/loadkeys
/bin/more
/bin/mount
/bin/ps
/bin/setfont
/bin/setserial
/bin/sort
/bin/tempfile
/bin/true
/etc/adjtime
/etc/conf.d/clock
/etc/conf.d/cryptfs
/etc/conf.d/hostname
/etc/conf.d/keymaps
/etc/csh.env
/etc/env.d/01hostname
/etc/fstab
/etc/hostname
/etc/init.d/bootmisc
/etc/init.d/checkfs
/etc/init.d/checkroot
/etc/init.d/clock
/etc/init.d/consolefont
/etc/init.d/functions.sh
/etc/init.d/hostname
/etc/init.d/keymaps
/etc/init.d/localmount
/etc/init.d/modules
/etc/init.d/net.lo
/etc/init.d/rmnologin
/etc/init.d/serial
/etc/init.d/urandom
/etc/modprobe.conf
/etc/modprobe.devfs
/etc/modules.autoload.d/kernel-2.4
/etc/modules.conf
/etc/modules.devfs
/etc/mtab
/etc/profile.env
/etc/rc.conf
/etc/resolv.conf
/etc/serial.conf
/etc/sysctl.conf
/lib/ld-linux.so.2
/lib/libacl.so.1
/lib/libattr.so.1
/lib/libblkid.so.1
/lib/libcrypt.so.1
/lib/libdl.so.2
/lib/libncurses.so.5
/lib/libnsl.so.1
/lib/libproc-3.2.5.so
/lib/libuuid.so.1
/lib/libz.so.1
/lib/rcscripts/awk/functions.awk
/lib/rcscripts/awk/genenviron.awk
/lib/tls/libc.so.6
/lib/tls/libm.so.6
/lib/tls/libpthread.so.0
/lib/tls/librt.so.1
/sbin/depmod
/sbin/env-update.sh
/sbin/functions.sh
/sbin/generate-modprobe.conf
/sbin/hwclock
/sbin/insmod
/sbin/losetup
/sbin/modprobe
/sbin/modules-update
/sbin/reboot
/sbin/sulogin
/sbin/swapon
/sbin/sysctl
/tmp/debug
/usr/bin/[
/usr/bin/as
/usr/bin/awk
/usr/bin/cat
/usr/bin/chmod
/usr/bin/cp
/usr/bin/current
/usr/bin/cut
/usr/bin/date
/usr/bin/dd
/usr/bin/diff
/usr/bin/disable
/usr/bin/display
/usr/bin/dumpkeys
/usr/bin/echo
/usr/bin/enable
/usr/bin/false
/usr/bin/file
/usr/bin/find
/usr/bin/gawk
/usr/bin/gpg
/usr/bin/hostname
/usr/bin/install
/usr/bin/kbd_mode
/usr/bin/line
/usr/bin/logger
/usr/bin/mkdir
/usr/bin/mv
/usr/bin/printenv
/usr/bin/prune
/usr/bin/rm
/usr/bin/sed
/usr/bin/seq
/usr/bin/setfont
/usr/bin/sleep
/usr/bin/sort
/usr/bin/sync
/usr/bin/tail
/usr/bin/touch
/usr/bin/tr
/usr/bin/true
/usr/bin/uname
/usr/bin/utmp
/usr/bin/xargs
/usr/bin/yes
/usr/lib/binutils/i686-pc-linux-gnu/2.15.92.0.2/libbfd-2.15.92.0.2.so
/usr/lib/graphviz/libagraph.so.0
/usr/lib/graphviz/libcdt.so.0
/usr/lib/libICE.so.6
/usr/lib/libMagick.so.6
/usr/lib/libSM.so.6
/usr/lib/libX11.so.6
/usr/lib/libXext.so.6
/usr/lib/libXi.so.6
/usr/lib/libXt.so.6
/usr/lib/libbz2.so.1.0
/usr/lib/libcrypto.so.0.9.7
/usr/lib/libcups.so.2
/usr/lib/libfreetype.so.6
/usr/lib/libgdk-1.2.so.0
/usr/lib/libglib-1.2.so.0
/usr/lib/libgmodule-1.2.so.0
/usr/lib/libgthread-1.2.so.0
/usr/lib/libgtk-1.2.so.0
/usr/lib/libjpeg.so.62
/usr/lib/libltdl.so.3
/usr/lib/libmagic.so.1
/usr/lib/libssl.so.0.9.7
/usr/lib/libtiff.so.3
/usr/lib/libxmms.so.1
/usr/local/bin/n
/usr/src/linux/System.map
