# scripts/readahead._sbin_rc.list
/bin/bash
/bin/dmesg
/bin/egrep
/bin/gawk
/bin/grep
/bin/more
/bin/mount
/bin/tar
/etc/init.d/functions.sh
/etc/init.d/halt.sh
/etc/init.d/reboot.sh
/etc/init.d/shutdown.sh
/etc/issue
/etc/rc.conf
/lib/ld-linux.so.2
/lib/libacl.so.1
/lib/libattr.so.1
/lib/libblkid.so.1
/lib/libcom_err.so.2
/lib/libcrypt.so.1
/lib/libdl.so.2
/lib/libe2p.so.2
/lib/libext2fs.so.2
/lib/libncurses.so.5
/lib/libuuid.so.1
/lib/rcscripts/awk/cachedepends.awk
/lib/rcscripts/awk/functions.awk
/lib/rcscripts/awk/gendepends.awk
/lib/tls/libc.so.6
/lib/tls/libm.so.6
/lib/tls/libpthread.so.0
/lib/tls/librt.so.1
/lib/udev-state/devices.tar.bz2
/sbin/depscan.sh
/sbin/devfsd
/sbin/functions.sh
/sbin/hotplug
/sbin/mke2fs
/sbin/rc
/sbin/reboot
/sbin/sulogin
/sbin/swapon
/sbin/telinit
/sbin/udev
/sbin/udevsend
/sbin/udevstart
/usr/bin/[
/usr/bin/as
/usr/bin/awk
/usr/bin/chgrp
/usr/bin/chmod
/usr/bin/clear
/usr/bin/cp
/usr/bin/dd
/usr/bin/echo
/usr/bin/hostname
/usr/bin/killall
/usr/bin/ln
/usr/bin/ls
/usr/bin/mkdir
/usr/bin/mv
/usr/bin/rm
/usr/bin/test
/usr/bin/touch
/usr/bin/utmp
/usr/bin/yes
/usr/lib/binutils/i686-pc-linux-gnu/2.15.92.0.2/libbfd-2.15.92.0.2.so
/usr/share/baselayout/issue.devfix
